# deploy/ralph-dispatcher@.service
[Unit]
Description=Ralph Loop — Dispatcher
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=%i
WorkingDirectory=%h/workspace/ClaudeXingCode
EnvironmentFile=%h/workspace/ClaudeXingCode/agent/.env
ExecStart=%h/workspace/ClaudeXingCode/.venv/bin/python3 %h/workspace/ClaudeXingCode/agent/dispatcher/dispatcher.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ralph-dispatcher

[Install]
WantedBy=multi-user.target
