# deploy/ralph-dispatcher.service
[Unit]
Description=Ralph Loop — Dispatcher
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/ClaudeXingCode
EnvironmentFile=/root/ClaudeXingCode/agent/.env
ExecStart=/root/ClaudeXingCode/.venv/bin/python3 /root/ClaudeXingCode/agent/dispatcher/dispatcher.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ralph-dispatcher

[Install]
WantedBy=multi-user.target
