# deploy/ralph-web.service
[Unit]
Description=Ralph Loop — Web Manager (Flask UI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/ClaudeXingCode
EnvironmentFile=/root/ClaudeXingCode/agent/.env
ExecStart=/root/ClaudeXingCode/.venv/bin/python3 /root/ClaudeXingCode/agent/web/web_manager.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ralph-web

[Install]
WantedBy=multi-user.target
