# deploy/ralph-web@.service
[Unit]
Description=Ralph Loop — Web Manager (Flask UI)
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/workspace/ClaudeXingCode
EnvironmentFile=%h/workspace/ClaudeXingCode/agent/.env
ExecStart=%h/workspace/ClaudeXingCode/.venv/bin/python3 %h/workspace/ClaudeXingCode/agent/web/web_manager.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ralph-web

[Install]
WantedBy=multi-user.target
